# payload/Library/LaunchDaemons/com.syncpool.app-downloader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.syncpool.app-downloader</string>
	<key>RunAtLoad</key>
	<true/>
	<key>ProgramArguments</key>
    <array>
        <string>/Library/Application Support/com.syncpool.app-downloader/app-downloader.sh</string>
    </array>
	<key>StandardErrorPath</key>
    <string>/Library/Application Support/com.syncpool.app-downloader/app-downloader.stderr</string>
	<key>StandardOutPath</key>
    <string>/Library/Application Support/com.syncpool.app-downloader/app-downloader.stdout</string>
</dict>
</plist>
